# deploy/aws/agenthook.service
[Unit]
Description=agenthook.store API service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
EnvironmentFile=/opt/agenthook/.env
ExecStart=/usr/local/bin/agenthook-server
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/agenthook
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
